# scraper/supervisord.conf
[supervisord]
nodaemon=true


[program:celeryworker]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=celery -A worker worker 


[program:celerybeat]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=celery -A worker beat


[program:notebook]
command=jupyter notebook --port=8888 --no-browser --ip=0.0.0.0 --allow-root


[program:celeryflower]
command = /usr/src/app/scraper/start-flower.sh
